# ansible/roles/soyoustart/files/fuck-you-iptables.service
[Unit]
Description=Fuck you iptables
After=network.target

[Service]
ExecStart=/sbin/iptables -F
Type=oneshot

[Install]
WantedBy=basic.target
