# systemd/coreos/fleet.socket
[Socket]
ListenStream=/var/run/fleet.sock

[Install]
WantedBy=multi-user.target
