# systemd/docker/dockerclean.service
[Unit]
Description=Clean orphaned containers and images

[Service]
Type=simple
ExecStart=-/usr/local/bin/cleanup-docker
