# systemd/pi/sync.service
[Unit]
Description=Sync filesystem on shutdown

[Service]
Type=oneshot
ExecStart=sync

[Install]
WantedBy=halt.target
